# ui-finger/rwho-finger@.service
[Unit]

[Service]
ExecStart=/usr/sbin/in.rwho-fingerd
StandardInput=socket
