# agent-linux/rwho-agent.service
[Unit]

[Service]
Type=forking
ExecStart=/usr/sbin/rwho-agent --daemon

[Install]
WantedBy=multi-user.target
